# conf.d/upstart-statsd.conf
# Statsd         
# ToDO: privileges - run as www-data or other user 
# http://superuser.com/questions/213416/running-upstart-jobs-as-unprivileged-users                                                                                        


description     "Statsd"

start on runlevel [2345]
stop on runlevel [!2345]


pre-start script
    echo 'start' >> /tmp/rhap.log
    echo '*** Start - statsd service' | logger -ist 'devweb-statsd'	 
end script

post-stop script
    echo 'stop' >> /tmp/rhap.log
    echo '*** Stop statsd service' | logger -ist 'devweb-statsd'
end script


script
 cd <<statsdhome>>
 /usr/bin/node stats.js dConfig.js 2>&1 | logger -ist 'devweb-statsd'
end script

respawn
respawn limit 3 20
